# deploy/nwz-web-frontend.service
[Unit]
Description=NWZ-Bot Web Frontend (Next.js)
After=network.target nwz-web-api.service

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/app/web/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=BACKEND_URL=http://127.0.0.1:8000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
